# apexpulsefxbot.service
[Unit]
Description=ApexPulseFXBot Service
After=network.target

[Service]
User=novak1873
WorkingDirectory=/path/to/apexpulsefxbot
ExecStart=/usr/bin/python3 main.py
Restart=always

[Install]
WantedBy=multi-user.target